# docker/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stdout_logfile          = /var/log/php-fpm-stdout.log
stdout_logfile_maxbytes = 0
stdout_events_enabled   = true
stderr_logfile          = /var/log/php-fpm-stderr.log
stderr_logfile_maxbytes = 0
stderr_events_enabled   = true

